# scriptmodules/supplementary/steamcontroller/steamcontroller-desktop.service
[Unit]
Description=steamcontroller-desktop
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/python /opt/archrgs/supplementary/steamcontroller/bin/sc-desktop.py start

[Install]
WantedBy=multi-user.target
